# regulatrix/files/regulatrix.init
#!/bin/sh /etc/rc.common
#
# Init script for regulatrix

START=95
PROG=/usr/sbin/regulatrix

start() {
	$PROG start
}

stop() {
	$PROG stop
}
